# horus/deploy/systemd/horus-paper.service
[Unit]
Description=Horus Paper Trading Runtime
After=network.target postgresql.service
Wants=network.target

[Service]
Type=simple
WorkingDirectory=%h/horus
EnvironmentFile=%h/horus/.env
ExecStart=%h/horus/.venv/bin/python -m horus.runtime.run_paper
Restart=always
RestartSec=3
TimeoutStopSec=10
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target
