# horus/deploy/systemd/horus-nightly.timer
[Unit]
Description=Run Horus nightly at 00:10 UTC

[Timer]
OnCalendar=*-*-* 00:10:00 UTC
Persistent=true

[Install]
WantedBy=default.target
