# horus/deploy/systemd/horus-nightly.service
[Unit]
Description=Horus Nightly Calibration + Gates + Reports

[Service]
Type=oneshot
WorkingDirectory=%h/horus
EnvironmentFile=%h/horus/.env
ExecStart=%h/horus/.venv/bin/python -m horus.ops.run_nightly
